# systemd/user/caelestia-visual-index.timer
[Unit]
Description=Continue visual search indexing in safe AC-powered batches

[Timer]
OnBootSec=1h
OnUnitActiveSec=30min
RandomizedDelaySec=10min
Persistent=true
Unit=caelestia-visual-index.service

[Install]
WantedBy=timers.target
