# systemd/user/caelestia-semantic-backfill.timer
[Unit]
Description=Finish the initial semantic index in small AC-powered batches

[Timer]
OnBootSec=15min
OnUnitActiveSec=30min
RandomizedDelaySec=15min
Persistent=true
Unit=caelestia-semantic-index.service

[Install]
WantedBy=timers.target
