# systemd/user/caelestia-clipboard.service
[Unit]
Description=Caelestia event-driven clipboard and Continuity bridge
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/home/yash/.local/bin/caelestia-clipboard serve
Restart=on-failure
RestartSec=2
Environment=PYTHONUNBUFFERED=1
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
