# systemd/user/caelestia-kdeconnect.service
[Unit]
Description=KDE Connect transport for Caelestia Continuity
After=graphical-session.target network.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/kdeconnectd
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
