# systemd/user/caelestia-cache-maintenance.timer
[Unit]
Description=Run low-impact cache maintenance weekly

[Timer]
OnCalendar=Sun *-*-* 04:00:00
RandomizedDelaySec=2h
Persistent=true

[Install]
WantedBy=timers.target
